# config/secure-storage.service
[Unit]
Description=Secure Storage Overlay (LUKS Loop)
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
# Assumes /home/debian/secure/storage.img exists and is LUKS formatted
# You must manually create it:
# dd if=/dev/zero of=/home/debian/secure/storage.img bs=1M count=1024
# cryptsetup luksFormat /home/debian/secure/storage.img
#
# This unit assumes a keyfile exists or requires manual interaction (which fails in systemd).
# For automation, use a keyfile protected by root:
# cryptsetup luksOpen /home/debian/secure/storage.img secure_storage --key-file /root/secure.key
ExecStart=/bin/true
# Placeholder: Real implementation depends on key management strategy.

[Install]
WantedBy=multi-user.target
